# contrib/window-bot.service
[Unit]
Description=Window Bot weather notifier
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/window-bot
EnvironmentFile=-/etc/window-bot/window-bot.env
Environment=PYTHONPATH=/opt/window-bot/src
ExecStart=/usr/bin/env python3 -m window_bot run --config /etc/window-bot/config.toml
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
